# system_files/shared/usr/lib/systemd/user/cosmic-keyring-env.service
[Unit]
Description=Update D-Bus activation environment for COSMIC keyring
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/dbus-update-activation-environment --systemd WAYLAND_DISPLAY DISPLAY
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
